# deploy/2026-05-07/systemd/auto-timezone.service
[Unit]
Description=Auto-detect and set timezone from IP geolocation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto-timezone.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
