# deploy/2026-05-07/systemd/haproxy-heal.timer
[Unit]
Description=Check HAProxy env vars every 10 seconds

[Timer]
OnBootSec=15
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
