# deploy/2026-05-07/systemd/jupyter-iptables-fix.service
[Unit]
Description=Block Docker bridge IPs in WebRTC ICE candidates
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/jupyter-iptables-fix.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
