# deploy/2026-05-07/systemd/jupyter-env-backup.service
[Unit]
Description=SecureProtect Hub .env Auto-Backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-hub-env.sh
