# deploy/2026-05-07/systemd/jupyter-docker-hygiene.timer
[Unit]
Description=Run Docker disk hygiene daily

[Timer]
OnBootSec=30min
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
