# deploy/2026-05-07/systemd/jupyter-disk-cleanup.service
[Unit]
Description=SecureProtect Disk Cleanup — oldest recordings first at 90%% disk
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/jupyter-disk-cleanup/disk-cleanup.sh
StandardOutput=journal
StandardError=journal
