# deploy/2026-05-07/systemd/jupyter-metrics.service
[Unit]
Description=Jupyter custom metrics emitter
[Service]
Type=oneshot
ExecStart=/usr/local/bin/jupyter-metrics-collector.sh
