# deploy/2026-05-07/systemd/jupyter-ota-lockdown.service
[Unit]
Description=SecureProtect OTA Lockdown - Protect baked files from OTA overwrites
After=docker.service
Before=hub-manager.service jupyter-hub-controller.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/ota-lockdown.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
