# deploy/2026-05-07/systemd/nvme-smart-exporter.service
[Unit]
Description=NVMe SMART exporter (textfile collector for node-exporter)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nvme-smart-exporter.sh
