# deploy/2026-05-07/systemd/nvme-smart-exporter.timer
[Unit]
Description=Run nvme-smart-exporter every hour

[Timer]
OnBootSec=30
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
