# deploy/2026-05-07/systemd/secureprotect-ip-reporter.timer
[Unit]
Description=Report Public IP every 5 minutes
[Timer]
OnBootSec=30s
OnUnitActiveSec=5min
[Install]
WantedBy=timers.target
