# deploy/2026-05-07/systemd/turn-credential-exporter.service
[Unit]
Description=TURN credential freshness exporter for Argus
After=postgres.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/turn-credential-exporter.sh
