# deploy/2026-05-07/systemd/turn-credential-exporter.timer
[Unit]
Description=Run turn-credential-exporter every 5 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
