# deploy/2026-05-07/systemd/jupyter-ble.service.d/override.conf
[Service]
Environment="JUPYTER_HUB_API_HOST=https://api.dev.jupyter.com.au"
